# nvim/com.llamacpp.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.llamacpp.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/llama-server</string>
        <string>-hf</string>
        <string>ggml-org/Qwen2.5-Coder-1.5B-Q8_0-GGUF</string>
        <string>--port</string>
        <string>8012</string>
        <string>-ngl</string>
        <string>99</string>
        <string>-fa</string>
        <string>auto</string>
        <string>-ub</string>
        <string>1024</string>
        <string>-b</string>
        <string>1024</string>
        <string>--ctx-size</string>
        <string>0</string>
        <string>--cache-reuse</string>
        <string>256</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/llamacpp-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/llamacpp-server.err</string>
</dict>
</plist>
